# misc/uwsgi/vis.conf
description "Vis uWSGI Instance"

start on runlevel [2345]
stop on runlevel [!2345]

setuid vis
setgid www-data

env PATH=/home/vis/vis/venv/bin
chdir /home/vis/vis
exec uwsgi --ini vis.ini
